# config/supervisor-openvscode.conf
[program:openvscode-server]
command=/usr/local/bin/openvscode-supervised
directory=/workspace
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=3
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/openvscode-server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/openvscode-server.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
