# config/supervisor-neat-insight.conf
[program:neat-insight]
# Launched via a wrapper that sources ~/.devkit-sync.rc so the service inherits
# the current CONTAINER_HOST_IP after a DevKit re-point. devkit.sh restarts this
# program on every (re-)point. Insight stays env-agnostic; the rc knowledge is
# confined to the wrapper.
command=/usr/local/bin/neat-insight-supervised
directory=/workspace
autostart=true
autorestart=true
startsecs=3
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/neat-insight.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/neat-insight.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
