# config/supervisor-neat-insight.conf
[program:neat-insight]
command=/opt/neat-insight/venv/bin/neat-insight --port %(ENV_NEAT_INSIGHT_PORT)s
directory=/workspace
autostart=true
autorestart=true
startsecs=3
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/neat-insight.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/neat-insight.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
